# deployment/jeeves-stats-web.service
[Unit]
Description=Jeeves Stats Web Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/zote/bots/jeeves
Environment="PATH=%h/.pyenv/shims:%h/.pyenv/bin:/usr/bin:/bin"

# Start the stats web server
ExecStart=%h/.pyenv/shims/python %h/bots/jeeves/stats_web.py --host 127.0.0.1 --port 8081

# Restart on failure
Restart=on-failure
RestartSec=5s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jeeves-stats-web

[Install]
WantedBy=default.target
